# etc/systemd/system/amedas.timer
[Unit]
Description = Amedas timer

[Timer]
OnCalendar = *-*-* *:0/10:10
AccuracySec = 1s

[Install]
WantedBy = timers.target
