# etc/systemd/system/bme280.timer
[Unit]
Description = BME280 logging timer

[Timer]
OnCalendar = *-*-* *:*:20
AccuracySec = 1s

[Install]
WantedBy = timers.target
