# etc/systemd/system/mhz19.timer
[Unit]
Description = MH-Z19 logging timer

[Timer]
OnCalendar = *-*-* *:*:40
AccuracySec = 1s

[Install]
WantedBy = timers.target
