# etc/systemd/system/speedtest.timer
[Unit]
Description = SpeedTest timer

[Timer]
OnCalendar = *-*-* *:05/10:00
AccuracySec = 1s

[Install]
WantedBy = timers.target
